# telegram-research-bot.service
[Unit]
Description=Telegram Research Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pantherstock
WorkingDirectory=/home/pantherstock/telegram-research-bot
EnvironmentFile=/home/pantherstock/telegram-research-bot/.env
ExecStart=/home/pantherstock/telegram-research-bot/venv/bin/python main.py
Restart=on-failure
RestartSec=5s
StandardOutput=append:/home/pantherstock/telegram-research-bot/logs/bot.log
StandardError=append:/home/pantherstock/telegram-research-bot/logs/bot.log

[Install]
WantedBy=multi-user.target
